# GZ服务器/qwen3-fallback.service
[Unit]
Description=Qwen3.5-2B Fallback LLM (CPU优化版 16K MTB 关闭思考)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/mnemosyne
ExecStartPre=/bin/sleep 3
ExecStart=/opt/mnemosyne/bin/llama-server \
  --port 11437 \
  --host 127.0.0.1 \
  --model /opt/mnemosyne/models/Qwen3.5-2B-Q4_K_M.gguf \
  --n-gpu-layers 0 \
  --threads 4 \
  --threads-batch 4 \
  --ctx-size 16384 \
  --batch-size 1024 \
  --ubatch-size 512 \
  --cont-batching \
  --mlock \
  --no-mmap \
  --cache-type-k q8_0 \
  --cache-type-v q8_0 \
  --reasoning off \
  --temp 0.2 \
  --top-p 0.6 \
  --top-k 10 \
  --min-p 0.03 \
  --presence-penalty 0.3 \
  --parallel 1 \
  --no-warmup
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
