# scripts/mnemosyne.service
[Unit]
Description=Mnemosyne Memory Engine v2.0
After=network.target postgresql@16-main.service
Requires=postgresql@16-main.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/mnemosyne
EnvironmentFile=/opt/mnemosyne/.env
Environment="PATH=/opt/mnemosyne/venv/bin:/usr/bin"
ExecStart=/opt/mnemosyne/venv/bin/python3 /opt/mnemosyne/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
